# conf/supervisor.conf
# -*- conf -*-

[program:apts]
command=/srv/sites/apts/conf/startup
user=ebower
autostart=true
autorestart=true
stderr_logfile=/srv/log/apts/supervisor_err.log
