# conf/supervisor.conf
# -*- conf -*-

[program:dcmdb]
command=/srv/sites/dcmdb/conf/startup
user=ebower
autostart=true
autorestart=true
stderr_logfile=/srv/log/dcmdb/supervisor_err.log
